# docker/php-fpm/laravel-queue-supervisor.conf
[supervisord]
nodaemon=true

[program:laravel-queue]
command = php /app/artisan queue:listen --sleep 5
user = root
autostart = true